# systemd/hermes-desktop-headless.service
[Unit]
Description=Hermes Desktop headless (Xvfb + noVNC)
Documentation=https://github.com/PabloTheThinker/hermes-desktop-headless
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=HD_DISPLAY=99
Environment=HD_BIND=127.0.0.1
Environment=HD_VNC_PORT=5901
Environment=HD_NOVNC_PORT=6080
ExecStart=%h/.local/bin/hermes-desktop-headless start
ExecStop=%h/.local/bin/hermes-desktop-headless stop
TimeoutStartSec=120

[Install]
WantedBy=default.target
